# build/projects/AX620Q_nand_arm32_k419_sipeed_nanokvm/rootfs/etc/systemd/system/ifplugd@.service
[Unit]
Description=ifplugd for network interface %I
Requires=ifup@%i.service
After=network.target ifup@%i.service

[Service]
Type=simple
ExecStart=/usr/sbin/ifplugd -n -i %I -f -u0 -d5 -I
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
